# deploy_tools/gunicorn.service.j2
[Unit]
Description=Gunicorn server for {{ ansible_host }}

[Service]
Restart=on-failure
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}/sites/{{ ansible_host }}
EnvironmentFile=/home/{{ ansible_user }}/sites/{{ ansible_host }}/.env

ExecStart=/home/{{ ansible_user }}/sites/{{ ansible_host }}/virtualenv/bin/gunicorn \
    --bind unix:/tmp/{{ ansible_host }}.socket \
    --access-logfile ../access.log \
    --error-logfile ../error.log \
    superlists.wsgi

[Install]
WantedBy=multi-user.target